# deploy/systemd/jip-india.service
[Unit]
Description=JIP Horizon India API
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/jip-india
EnvironmentFile=/opt/jip-india/.env
ExecStart=/opt/jip-india/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8001 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
